# service/exaproxy/run
#!/bin/sh
exec nice --4 env PYTHONPATH=lib sbin/exaproxy 2>&1
